# .config/systemd/user/hyprlock.service
[Service]
EnvironmentFile=%t/hyprlock.env
ExecSearchPath=/run/wrappers/bin:/run/current-system/sw/bin
ExecStart=exec-after hyprlock & sudo chvt 3
