# .config/systemd/user/lock.service
[Unit]
UpheldBy=lock.target
Before=lock.target
OnSuccess=unlock.target

[Service]
ExecSearchPath=/run/current-system/sw/bin
ExecStart=hyprlock
KillSignal=SIGUSR1
Restart=on-failure
RestartSec=0

[Install]
WantedBy=lock.target
