# .config/systemd/user/mr.service
[Service]
ExecSearchPath=/run/current-system/sw/bin
ExecStart=mr up
