# config/svc/ux/lock.service
[Unit]
BindsTo=lock.target
Before=lock.target
OnSuccess=unlock.target

[Service]
ExecSearchPath=/run/current-system/sw/bin
ExecStart=hyprlock
KillSignal=SIGUSR1

[Install]
UpheldBy=lock.target
