# config/svc/ux/shell@.service
[Unit]
Requisite=shell.service
After=shell.service

[Service]
Type=oneshot
ExecStart=/bin/shell %I
