# config/svc/ux/sleep-reload.service
[Unit]
Before=lock.service
Requisite=hyprland.service

[Service]
Type=oneshot
ExecStart=/bin/hyprctl reload

[Install]
WantedBy=sleep.target
