# config/svc/ux/wake.target
[Unit]
Conflicts=sleep.target
Before=sleep.target

[Install]
WantedBy=default.target
