# config/svc/msft/work.service
[Unit]
Conflicts=sleep.target
Before=sleep.target
Wants=intune.service
PropagatesStopTo=intune.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/bin/hyprctl eval 'view.disable "eDP-1"'
