# config/svc/misc/term-sleep.service
[Unit]
BindsTo=sleep.target
Before=sleep.target

[Service]
Type=oneshot
ExecStart=/bin/sudo chvt 2

[Install]
WantedBy=sleep.target
