# config/svc/blade/hyprland.service.d/overrides.conf
[Unit]
Requires=dev-user-nvidia.device dev-user-intel.device
After=dev-user-nvidia.device dev-user-intel.device
